# extracted/_dlink.BIN.extracted/squashfs-root/etc/inittab
::sysinit:/etc/rc.d/rcS
ttyS0::askfirst:/bin/ash
#::shutdown:/etc/ath/killVAP all
#::shutdown:/bin/umount -a
